# deploy/max-university-backend.service
# Пример unit для systemd: бэкенд max-university.
# Установка:
#   sudo cp deploy/max-university-backend.service /etc/systemd/system/
#   sudo sed -i "s|/path/to/repo|$(pwd)|g" /etc/systemd/system/max-university-backend.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable max-university-backend
#   sudo systemctl start max-university-backend
#
# Переменные EVENTS_API_URL и EVENTS_BOT_LINK можно задать в backend/.env.events
# (создаётся при деплое из секретов GitHub) или здесь через Environment=.
# Для ленты новостей: Environment="COLD_NEWS_FEED_URL=http://127.0.0.1:3001"
# Для PostgreSQL: создайте backend/.env.database с DATABASE_URL=postgresql://...
# Токен бота MAX: создайте backend/.env.bot с одной строкой MAX_BOT_TOKEN=ваш_токен

[Unit]
Description=Max University FastAPI Backend
After=network.target

[Service]
Type=simple
User=www
Group=www
WorkingDirectory=/path/to/repo/backend
Environment="PATH=/path/to/repo/venv/bin"
EnvironmentFile=/path/to/repo/backend/.env.events
EnvironmentFile=-/path/to/repo/backend/.env.database
EnvironmentFile=-/path/to/repo/backend/.env.bot
Environment="COLD_NEWS_FEED_URL=http://127.0.0.1:3001"
ExecStart=/path/to/repo/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
